# etc_systemd/auto_spool_gate.service
[Unit]
Description=Auto Spool Gate NFC/MMU Service
After=network.target

[Service]
Type=simple
User=mks
WorkingDirectory=/home/mks
ExecStart=/usr/bin/python3 /home/mks/auto_spool_gate.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
